# scripts/dog-voice.service
[Unit]
Description=Voice Interface V2 Dog Receiver
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/unitree/Voice-Interface-V2/Voice-Interface-V2
ExecStartPre=/bin/sh /home/unitree/Voice-Interface-V2/Voice-Interface-V2/scripts/dog_set_wired_ip.sh
ExecStart=/usr/bin/python3 -u /home/unitree/Voice-Interface-V2/Voice-Interface-V2/dog/receiver.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
